# rootfs/etc/supervisor.d/dashing.ini

[program:dashing]
;user = www-data
directory      = /opt/dashing
command        = run.sh
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true

